# infra/systemd/skill-hub-api.service
[Unit]
Description=Skill Hub API
After=network.target postgresql.service redis.service
Wants=network.target

[Service]
Type=simple
User=skillhub
Group=staff
WorkingDirectory=/Users/xusheng/Documents/project/cjhx-skill-hub/apps/api-server
EnvironmentFile=/Users/xusheng/Documents/project/cjhx-skill-hub/.env.production
ExecStart=/Users/xusheng/Documents/project/cjhx-skill-hub/apps/api-server/.venv/bin/gunicorn app.main:app --workers 4 --worker-class uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --timeout 120 --access-logfile - --error-logfile -
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
